# batch/crontab
# Run batch jobs every day at 02:00 UTC.
# D must run first — stream job 3 reads vessel profiles from MongoDB.
# A, B, C have no mutual dependencies and run after D.
0  2 * * * /opt/spark-jobs/batch/run_batch_jobs.sh >> /var/log/batch_jobs.log 2>&1
